# ops/fastify-api.service.example
[Unit]
Description=Cabras Chat
After=network.target

[Service]
Type=simple
User=giovanni
WorkingDirectory=/srv/apps/cabras-chat
Environment=HOST=127.0.0.1
Environment=PORT=3000
ExecStart=/usr/bin/node /srv/apps/cabras-chat/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
